# nginx/backup/eplatform-django.service
[Unit]
Description=E_Platform Django Service
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
User=http
Group=http
WorkingDirectory=/mnt/4BC9EBFC28ECC8F5/Others/Projects/E_Platform
ExecStart=/mnt/4BC9EBFC28ECC8F5/Others/Projects/E_Platform/.venv/bin/gunicorn \
    --workers 3 \
    --bind 127.0.0.1:8000 \
    --timeout 120 \
    --access-logfile /var/log/eplatform/gunicorn-access.log \
    --error-logfile /var/log/eplatform/gunicorn-error.log \
    E_Platform.wsgi:application

# Environment variables
Environment=DB_ENGINE=postgresql
Environment=DB_NAME=e_platform_db
Environment=DB_USER=zero
Environment=DB_PASSWORD=82821931003
Environment=DB_HOST=localhost
Environment=DB_PORT=5432

# Restart on failure
Restart=on-failure
RestartSec=5s

# Limit resources
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
